# rgc.service
[Unit]
Description=Remote GPIO Control Server
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python rgc-server.py
Restart=always
WorkingDirectory=/home/pi
RestartSec=3

[Install]
WantedBy=multi-user.target

